# deploy/systemd/pipeworks-character-forge.service
[Unit]
Description=PipeWorks Character Forge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pipeworks
Group=pipeworks
WorkingDirectory=/srv/work/pipeworks/repos/pipeworks-character-forge
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/pipeworks/character-forge/character-forge.env
ExecStart=/srv/work/pipeworks/venvs/pw-character-forge/bin/pw-forge
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
